# deploy/systemd/metube-telegram-bot.service
[Unit]
Description=MeTube Telegram Proxy Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/metube-telegram-proxy-bot
EnvironmentFile=/opt/metube-telegram-proxy-bot/.env
ExecStart=/usr/bin/python3 -m bot.main
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
